# services/watchdog/watchdog.service
[Unit]
Description=AlphaDivision LLM Watchdog
After=network.target ollama.service docker.service
Wants=ollama.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/alphadivision/services/watchdog
EnvironmentFile=/opt/alphadivision/.env
Environment=CONFIG_FILE=/opt/alphadivision/config.toml
ExecStart=/usr/bin/python3 /opt/alphadivision/services/watchdog/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
